# ops/systemd/bus-route-health-guardian.service
[Unit]
Description=Bus Route Demo Host Health Guardian
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/dev/bus-route-query-system
ExecStart=/root/dev/bus-route-query-system/scripts/host-health-guardian.sh
